# pkg/mdaffin/mdaffin-desktop.install
post_install() {
    post_upgrade
    systemctl enable --now sddm
    systemctl enable --now connman
    systemctl enable --now avahi-daemon
}

post_upgrade() {
    cp /etc/i3/mdaffin-desktop-config /etc/i3/config 
    cp /etc/xdg/termite/mdaffin-desktop.config /etc/xdg/termite/config
}
