# pkg/base/mdaffin-base.install
post_install() {
    systemctl enable --now systemd-timesyncd
    ln -sf /usr/share/zoneinfo/Europe/London "$pkgdir/etc/localtime"
    locale-gen
}

post_upgrade() {
    locale-gen
}
