# pkg/base/mdaffin-base.install
post_install() {
    systemctl enable --now systemd-timesyncd
    ln -sf /usr/share/zoneinfo/Europe/London "/etc/localtime"
    sed 's/#en_GB/en_GB/' -i /etc/locale.gen
    locale-gen
}

post_upgrade() {
    sed 's/#en_GB/en_GB/' -i /etc/locale.gen
    locale-gen
}
