# pkg/mdaffin/devel/files/etc/profile.d/mdaffin-devel.sh
# Node
export PATH=$HOME/.local/share/npm/bin:$HOME/.yarn/bin:$PATH

# Go
export GOPATH=$HOME/projects/go
export PATH=$GOPATH/bin:$PATH

# Rust
export PATH=$HOME/.cargo/bin:$PATH
