# pkg/mdaffin/devel-languages/etc/profile.d/mdaffin-devel-go.sh
export GOPATH=$HOME/projects/go
export PATH=$GOPATH/bin:$PATH
