# pkg/desktop/mdaffin-desktop.install
post_install() {
    systemctl enable --now sddm
    systemctl enable --now connman
    systemctl enable --now avahi-daemon
}
